# code-server/supervisord.conf
[supervisord]
nodaemon=false
environment=HOME=/home/coder

[program:vs-code-2]
command=/usr/lib/code-server/bin/code-server --host 0.0.0.0 --port 8080 --auth=none
user=coder
directory=/home/coder

