# v1/matlab/containerd/supervisord.conf
[supervisord]
nodaemon=false
environment=HOME=/home/matlab

[program:matlab]
command=matlab-proxy-app
#user=matlab
directory=/home/matlab

