# tiger-vnc/supervisord.conf
[supervisord]
nodaemon=false
environment=HOME=/home/coder

[program:tigervnc]
command=vncserver :1 -SecurityTypes None
user=coder
directory=/home/coder



