# tiger-vnc/supervisord.conf
[supervisord]
nodaemon=false
environment=HOME=/home/coder

[program:tigervnc]
command=vncserver :1 -localhost no -passwd /home/coder/.vnc/passwd
user=coder
directory=/home/coder

