# pgadmin/supervisord.conf
[supervisord]
nodaemon=false
environment=HOME=/home/coder

[program:postgres]
command=/etc/init.d/postgresql start
user=coder
directory=/home/coder

