# v1/coder-training/web-server/supervisord.conf
[supervisord]
nodaemon=false
environment=HOME=/home/coder

[program:python-web-server]
command=python3 -m http.server 8888
user=coder
directory=/home/coder

