# airflow/supervisord/supervisord.conf
[supervisord]
nodaemon=false
environment=HOME=/home/coder

[program:airflow]
command=airflow standalone
user=coder
directory=/home/coder

